# infra/systemd/presence-snapshot-purge.timer
[Unit]
Description=Run Sheoak Tree presence snapshot purge daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
